# assets/config/supervisor/5waterfox.conf
[program:Waterfox]
command=bash -c 'waterfox --no-sandbox --test-type --disable-dev-shm-usage --disable-gpu --start-maximized'
priority=5
autorestart=true
stderr_logfile=/var/log/waterfox.err.log
